# crontab
0 */3 * * * /bin/bash /scripts/ustv.sh
30 */12 * * * /bin/bash /scripts/toonami.sh
55 * * * * /bin/bash /scripts/nhl.sh
